# software/system/etc/supervisor/conf.d/ei.conf
[program:ei]
command=/home/pi/eclipse/software/ei
autostart=true
autorestart=true
redirect_stderr=true
directory=/home/pi/eclipse/software
user=pi
