# software/system/etc/supervisor/conf.d/ei.conf
[program:ei]
command=/home/micah/src/ecstatic-epiphany/ei
autostart=true
autorestart=true
redirect_stderr=true
directory=/home/micah/src/ecstatic-epiphany
user=micah
